# fund_exit/com.spxstrat.fundexit.refresh.plist
<?xml version="1.0" encoding="UTF-8"?>
<!-- 老 Air 日度刷新基金清仓信号。安装见 fund_exit/README.md。
     部署目标: ~/Library/LaunchAgents/com.spxstrat.fundexit.refresh.plist (user macbook)
     路径为老 Air 绝对路径; 换机需同步调整。 -->
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key><string>com.spxstrat.fundexit.refresh</string>
  <key>ProgramArguments</key>
  <array>
    <string>/Users/macbook/SPX_strat/.venv_fund/bin/python</string>
    <string>/Users/macbook/SPX_strat/fund_exit/fund_exit_signals.py</string>
  </array>
  <key>WorkingDirectory</key><string>/Users/macbook/SPX_strat</string>
  <key>StartCalendarInterval</key>
  <array>
    <dict><key>Weekday</key><integer>1</integer><key>Hour</key><integer>14</integer><key>Minute</key><integer>0</integer></dict>
    <dict><key>Weekday</key><integer>2</integer><key>Hour</key><integer>14</integer><key>Minute</key><integer>0</integer></dict>
    <dict><key>Weekday</key><integer>3</integer><key>Hour</key><integer>14</integer><key>Minute</key><integer>0</integer></dict>
    <dict><key>Weekday</key><integer>4</integer><key>Hour</key><integer>14</integer><key>Minute</key><integer>0</integer></dict>
    <dict><key>Weekday</key><integer>5</integer><key>Hour</key><integer>14</integer><key>Minute</key><integer>0</integer></dict>
  </array>
  <key>StandardOutPath</key><string>/Users/macbook/SPX_strat/fund_exit/refresh.log</string>
  <key>StandardErrorPath</key><string>/Users/macbook/SPX_strat/fund_exit/refresh.log</string>
  <key>RunAtLoad</key><false/>
</dict>
</plist>
