# com.spxstrat.q085_s2bps.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN"
  "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>com.spxstrat.q085_s2bps</string>

  <key>ProgramArguments</key>
  <array>
    <string>/Users/macbook/SPX_strat/venv/bin/python</string>
    <string>-m</string>
    <string>notify.q085_s2bps_paper</string>
  </array>

  <key>WorkingDirectory</key>
  <string>/Users/macbook/SPX_strat</string>

  <!-- 16:50 ET Mon-Fri, after the 16:30 q041 chain snapshot -->
  <key>StartCalendarInterval</key>
  <array>
    <dict><key>Hour</key><integer>16</integer><key>Minute</key><integer>50</integer><key>Weekday</key><integer>1</integer></dict>
    <dict><key>Hour</key><integer>16</integer><key>Minute</key><integer>50</integer><key>Weekday</key><integer>2</integer></dict>
    <dict><key>Hour</key><integer>16</integer><key>Minute</key><integer>50</integer><key>Weekday</key><integer>3</integer></dict>
    <dict><key>Hour</key><integer>16</integer><key>Minute</key><integer>50</integer><key>Weekday</key><integer>4</integer></dict>
    <dict><key>Hour</key><integer>16</integer><key>Minute</key><integer>50</integer><key>Weekday</key><integer>5</integer></dict>
  </array>

  <key>RunAtLoad</key>
  <false/>

  <key>KeepAlive</key>
  <false/>

  <key>StandardOutPath</key>
  <string>/Users/macbook/SPX_strat/logs/q085_s2bps.log</string>

  <key>StandardErrorPath</key>
  <string>/Users/macbook/SPX_strat/logs/q085_s2bps.err.log</string>

  <key>EnvironmentVariables</key>
  <dict>
    <key>PATH</key>
    <string>/Users/macbook/SPX_strat/venv/bin:/usr/local/bin:/usr/bin:/bin</string>
    <key>HOME</key>
    <string>/Users/macbook</string>
  </dict>
</dict>
</plist>
